# Library/LaunchAgents/local.cleanup.pictures.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>local.cleanup.pictures</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/sh</string>
		<string>-c</string>
		<string>sh -c "$HOME/bin/cleanup-directory -d $HOME/Pictures -l 20 -v"</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardErrorPath</key>
	<string>/tmp/local.cleanup.pictures.stderr</string>
	<key>StandardOutPath</key>
	<string>/tmp/local.cleanup.pictures.stdout</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>3</integer>
		<key>Minute</key>
		<integer>32</integer>
	</dict>
</dict>
</plist>
