# Library/LaunchAgents/local.cleanup.trash.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>local.cleanup.trash</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/sh</string>
		<string>-c</string>
		<string>$HOMEBREW_PREFIX/sbin/tmpreaper -vv --showdeleted -M --delay -T 300 5d $HOME/.Trash</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardErrorPath</key>
	<string>/tmp/local.cleanup.trash.stderr</string>
	<key>StandardOutPath</key>
	<string>/tmp/local.cleanup.trash.stdout</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>3</integer>
		<key>Minute</key>
		<integer>40</integer>
	</dict>
</dict>
</plist>
